# src/server/budget-api.service
[Unit]
Description=Budget API backend service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/BudgetApp/BudgetApp/src/server
ExecStart=/usr/bin/node /home/pi/BudgetApp/BudgetApp/src/server/server.js
Restart=on-failure
RestartSec=5s
Environment=DB_TYPE=postgres
Environment=DB_CONNECTION_STRING=postgresql://admin:Password123@localhost:5432/homebudget
Environment=PORT=3000
User=pi
Group=pi
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
