# expo-app.service
[Unit]
Description=BudgetApp Expo web dev server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/BudgetApp/BudgetApp
ExecStart=/usr/bin/npm run web -- --host 0.0.0.0 --port 8081
Restart=on-failure
RestartSec=5s
Environment=EXPO_NO_BROWSER=1
User=pi
Group=pi
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
